# deploy/systemd/flash-note-health-watchdog.timer
[Unit]
Description=Run Flash Note health watchdog every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s
Persistent=true
Unit=flash-note-health-watchdog.service

[Install]
WantedBy=timers.target
