# deploy/systemd/flash-note-backup.service
[Unit]
Description=Shunnian Notes Backup Snapshot
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/flash-note
EnvironmentFile=-/root/flash-note/.env
ExecStart=/usr/bin/python3 /root/flash-note/scripts/backup_snapshot.py
